# rootfs/etc/services.d/pdns/run
#!/usr/bin/with-contenv sh

# Run pdns server
/bin/bash /xshok_gen_conf.sh
